# next-controller.service
[Unit]
Description=Next Controller
After=network.target

[Service]
Type=simple
WorkingDirectory=/mnt/nvme/next-controller
ExecStart=/mnt/nvme/next-controller/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
